# packages/server/halmp-server/systemd/black_shell.service
[Unit]
Description=HAL Media Player startup service
After=getty.target

[Service]
Type=idle
User=pi
ExecStart=sudo sh -c "TERM=linux setterm -foreground black -clear all >/dev/tty1"

[Install]
WantedBy=multi-user.target