# packages/server/halmp-server/systemd/halmp.service
[Unit]
Description=HAL Media Player startup service
After=network.target

[Service]
WorkingDirectory=/home/pi/HAL-MediaPlayer/packages/server
Restart=on-failure
User=pi
ExecStart=authbind --deep python -m halmp-server

[Install]
WantedBy=multi-user.target